# shifterboot.conf
description     "shifterboot startup config "
author "europa502<abhijit502@gmail.com>"
start on startup
task
exec /path/to/my/script.sh
